# modules/oftc/files/puppetmaster.cron.daily
#!/bin/sh

find /var/lib/puppet/reports -type f -mtime +7 | xargs -r rm
