# modules/oftc/files/puppetmaster.cron.daily
#!/bin/sh

find /var/lib/puppet/reports -type f -mtime +14 | xargs -r rm
